# zzzz.conf
[program:zzzz]
command=python3 /home/curator2/zzzz/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/zzzz/test.err.log
stdout_logfile=/home/curator2/zzzz/test.out.log